# .config/systemd/user/batmon.service
[Unit]
Description=Battery monitor

[Service]
ExecStart=/home/oli/.local/bin/batmonitor

