# .config/systemd/user/calsync.service
[Unit]
Description=Calendar Sync

[Service]
ExecStart=/home/oli/.local/bin/synccalendar
